# private_dot_config/systemd/user/backup.timer
[Unit]
Description=Run backup.service every week

[Timer]
OnBootSec=15min
OnUnitActiveSec=1w
Persistent=true

[Install]
WantedBy=timers.target
