# private_dot_config/systemd/user/backup.service
[Unit]
Description=Backup with script ~/.local/bin/backup

[Service]
Type=oneshot
ExecStart=%h/.local/bin/backup
