# private_dot_config/systemd/user/get-kerberos-ticket.service
[Unit]
Description=Get a Kerberos ticket

[Service]
Type=oneshot
ExecStart=kinit -k -t %h/.cache/krb5.keytab.ecorcgae nfs/ecorcgae@EINFRA
